# crontabs
*/1 * * * * /bin/main.sh
@reboot sleep 600 && /bin/startscript.sh
